# infra/aws/ray-reconnect.service
# scripts/ray-reconnect.service
[Unit]
Description=Ray Worker Reconnect
After=network-online.target
Wants=network-online.target
[Service]
Type=simple
ExecStart=/usr/local/bin/worker_reconnect.sh
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target
